# pkg/conf/backup/etcd/backup-etcd.timer
[Unit]
Description=Backup Etcd Every Hour

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
AccuracySec=1h

[Install]
WantedBy=timers.target
