# pkg/conf/backup/etcd/backup-etcd.service
[Unit]
Description=Backup Etcd Every Hour

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/backup-etcd.sh
