# pkg/conf/backup/etcd/backup-etcd.timer
[Unit]
Description=Backup Etcd Every Hour

[Timer]
OnCalendar=*-*-* 02:00:00
Unit=backup-etcd.service

[Install]
WantedBy=timers.target
