# pkg/conf/backup/etcd/backup-etcd.service
[Unit]
Description=My Task Service

[Service]
Type=oneshot
ExecStart=/path/to/your/command
